# .config/systemd/user/pass-secrets.service.d/override.conf
[Service]
Environment="PASSWORD_STORE_DIR=/home/zneix/.local/share/pass-secrets"
